# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log

[program:MariaDB]
command=/etc/init.d/mysql.server start

